# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:noit]
command=/usr/local/sbin/noitd -c noitconf/noit.conf

[program:stratcon]
command=/usr/local/sbin/stratcond -c noitconf/stratcon.conf